# supervisor/alipaycode.conf
[program:alipaycode]
command=php /var/www/html/plugins/alipaycode/server.php %(ENV_ALIPAYCODE_CHANNEL)s
directory=/var/www/html/plugins/alipaycode
user=www-data
autostart=true
autorestart=true
startretries=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
